# conf/goapp/supervisord.d/supervisor.ini
[program:deskmgr]
command=/home/working/www/deskmgr/deskmgr --config=/home/working/www/deskmgr/config/config.toml
directory=/home/working/www/deskmgr/
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; whether/when to restart (default: unexpected)
startsecs=1                   ; number of secs prog must stay running (def. 1)
startretries=9999                ; max # of serial start failures (default 3)
user=root                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/home/working/data/logs/sys/deskmgr.stdout.log        ; 
nodaemon=false ; no daemon

[program:myapi]
command=/home/working/www/myapi/myapi -c config/config.yaml
directory=/home/working/www/myapi/
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; whether/when to restart (default: unexpected)
startsecs=1                   ; number of secs prog must stay running (def. 1)
startretries=9999                ; max # of serial start failures (default 3)
user=root                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/home/working/data/logs/sys/myapi.stdout.log        ;
nodaemon=false ; no daemon
