# conf/gateway/supervisord.d/supervisor.ini
[program:openresty]
command=/usr/local/openresty/nginx/sbin/nginx -c conf/nginx.conf
directory=/usr/local/openresty
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; whether/when to restart (default: unexpected)
startsecs=1                   ; number of secs prog must stay running (def. 1)
startretries=9999                ; max # of serial start failures (default 3)
user=root                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/home/working/data/logs/sys/openresty.stdout.log        ; 
nodaemon=false ; no daemon 前台运行。

